# server/s6/opendkim/run
#!/command/execlineb -P

# Route stderr to stdout
fdmove -c 2 1

/usr/sbin/opendkim -f
