# server/s6/opendkim/run
#!/command/execlineb -P

# Pipe stderr to stdout and add prefix
pipeline { fdmove -c 2 1 /usr/sbin/opendkim -f } sed --unbuffered "s/^/\[opendkim\]: /g"
