# server/s6/opendkim/run
#!/command/execlineb -P
/usr/sbin/opendkim -f
