# s6/dovecot/run
#!/command/execlineb -P

# Pipe stderr to stdout and add prefix
pipeline { fdmove -c 2 1 /usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf } sed -e "s/^/\[dovecot\]: /g"
