# dashboard/s6/dashboard/run
#!/command/execlineb -P

# Pipe stderr to stdout and add prefix
pipeline { fdmove -c 2 1 /usr/local/bin/shell2http -cgi GET:/dashboard "cat /opt/dashboard/index.html" POST:/dashboard "/opt/dashboard/server.sh" GET:/dashboard/logo "cat /icons/logo" } sed -e "s/^/\[dashboard\]: /g"
