# server/s6/dovecot/run
#!/command/execlineb -P
/usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
