# s6/opendmarc/run
#!/command/execlineb -P

# Pipe stderr to stdout and add prefix
pipeline { fdmove -c 2 1 /usr/sbin/opendmarc -f -c /etc/opendmarc.conf } sed -e "s/^/\[opendmarc\]: /g"
