# s6/spamassassin/run
#!/command/execlineb -P
/usr/sbin/spamd --create-prefs --max-children 5 --helper-home-dir
